# deploy/production/systemd/greenfish-celerybeat.service
# Celery beat scheduler for GreenFish periodic jobs:
# notification outbox dispatch, unpaid-order expiry, session cleanup,
# GDPR data-request processing, and personal-data retention.
[Unit]
Description=GreenFish celery beat scheduler
After=network.target redis-server.service

[Service]
Type=simple
User=deploy
Group=www-data
WorkingDirectory=/home/deploy/greenfish
Environment=DJANGO_SETTINGS_MODULE=config.settings.production
ExecStart=/home/deploy/greenfish/venv/bin/celery -A config beat \
    --loglevel=INFO \
    --schedule /home/deploy/greenfish/.celerybeat-schedule
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
