# deploy/production/systemd/greenfish-gunicorn.service
# Gunicorn service for GreenFish.
# Adjust User/Group and paths to match your server, then:
#   sudo systemctl enable --now greenfish-gunicorn.socket greenfish-gunicorn.service
[Unit]
Description=GreenFish gunicorn daemon
Requires=greenfish-gunicorn.socket
After=network.target

[Service]
Type=notify
User=deploy
Group=www-data
WorkingDirectory=/home/deploy/greenfish
Environment=DJANGO_SETTINGS_MODULE=config.settings.production
ExecStart=/home/deploy/greenfish/venv/bin/gunicorn \
    --workers 3 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 60 \
    --graceful-timeout 15 \
    --access-logfile - \
    --error-logfile - \
    --bind unix:/run/greenfish-gunicorn.sock \
    config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=20
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
